# chromeos-base/chromeos-visl/files/visl.conf
# Copyright 2023 The ChromiumOS Authors
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description    "Install VISL (Virtual Stateless V4L2 video decoder) module"
author         "chromium-os-dev@chromium.org"

# Ensure the module is loaded before Chrome starts.
start on starting ui
task
oom score never

script

  if [ -d /sys/bus/platform/drivers/visl ]; then
    # Do nothing if driver is builtin or module is already loaded.
    exit 0
  fi

  if ! output=$(modprobe visl 2>&1); then
    echo "${output}" | logger -t "${UPSTART_JOB}" -p ERR
    exit 1
  fi

end script
