# chromeos-base/chromeos-gce-l1-vkms/files/vkms.conf
# Copyright 2023 The ChromiumOS Authors
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description    "Install vkms module when no VGA"
author         "chromium-os-dev@chromium.org"

# Ensure the module is loaded before Chrome starts.
start on starting ui
task
oom score never

script

  case "$(lspci)" in
    # "Display controller" works on arm64.
    *"Display controller"*)
      logger -t "${UPSTART_JOB}" "Display controller is present; do nothing"
      exit 0
      ;;
    # "VGA compatible" works on x86_64.
    *"VGA compatible"*)
      logger -t "${UPSTART_JOB}" "VGA device is present; do nothing"
      exit 0
      ;;
  esac

  if lsmod | grep -q "vkms"; then
    logger -t "${UPSTART_JOB}" \
        "vkms module is already loaded; do nothing"
    exit 0
  fi

  modprobe vkms 2>&1 | logger -t "${UPSTART_JOB}" -p ERR

end script
